# jsi.conf
description "jsi"

start on (filesystem)
stop on runlevel [016]
respawn
console log

setuid alorente
env HOME=/home/alorente
chdir /home/alorente/my-first-blog

exec /usr/local/bin/npm start
